# build/docker/root/etc/services.d/audiobookshelf/finish
#!/bin/sh

if [ "$1" -eq 256 ]; then
  e=$((128 + $2))
else
  e="$1"
fi

echo "$e" > /run/s6-linux-init-container-results/exitcode