# build/docker/root/etc/services.d/audiobookshelf/run
#!/usr/bin/with-contenv bash


PORT=${PORT:-80}

exec \
    s6-setuidgid abc npm --prefix /app/audiobookshelf/ --cache /app/audiobookshelf/ start --port "${PORT}"